# src/php-cron-5.6/conf/crontab
* * * * * /bin/sh /var/www/html/cron.sh >> /var/www/html/var/log/magento.cron.log
